# systemd/airflow-scheduler.service
[Unit]
Description=Airflow scheduler daemon
After=network.target postgresql.service
Wants=postgresql.service

[Service]
User=airflow
Group=airflow
Type=simple
ExecStart=/bin/bash -c 'cd /home/airflow/airflow && make scheduler'
Restart=on-failure
RestartSec=5s
PrivateTmp=true

[Install]
WantedBy=multi-user.target

